# contrib/init/picacoind.conf
description "Picacoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env PICACOIND_BIN="/usr/bin/picacoind"
env PICACOIND_USER="picacoin"
env PICACOIND_GROUP="picacoin"
env PICACOIND_PIDDIR="/var/run/picacoind"
# upstart can't handle variables constructed with other variables
env PICACOIND_PIDFILE="/var/run/picacoind/picacoind.pid"
env PICACOIND_CONFIGFILE="/etc/picacoin/picacoin.conf"
env PICACOIND_DATADIR="/var/lib/picacoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # picacoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$PICACOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run picacoind."
        echo "The setting must appear in $PICACOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Picacoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$PICACOIND_PIDDIR"
    chmod 0755 "$PICACOIND_PIDDIR"
    chown $PICACOIND_USER:$PICACOIND_GROUP "$PICACOIND_PIDDIR"
    chown $PICACOIND_USER:$PICACOIND_GROUP "$PICACOIND_CONFIGFILE"
    chmod 0660 "$PICACOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$PICACOIND_PIDFILE" \
    --chuid $PICACOIND_USER:$PICACOIND_GROUP \
    --exec "$PICACOIND_BIN" \
    -- \
    -pid="$PICACOIND_PIDFILE" \
    -conf="$PICACOIND_CONFIGFILE" \
    -datadir="$PICACOIND_DATADIR" \
    -disablewallet \
    -daemon

